# profile.d/20-modules.ps1
# Various PowerShell modules
_Import-ModuleIfAvailable PSReadLine
_Import-ModuleIfAvailable AngleParse
