# profile.d/60-winget.ps1
# winget autocomplete (only if winget is installed)
if (Get-Command -Name winget -ErrorAction SilentlyContinue)
{
  Register-ArgumentCompleter -Native -CommandName winget -ScriptBlock {
    param($wordToComplete, $commandAst, $cursorPosition)
        [Console]::InputEncoding = [Console]::OutputEncoding = $OutputEncoding = [System.Text.Utf8Encoding]::new()
        $Local:word = $wordToComplete.Replace('"', '""')
        $Local:ast = $commandAst.ToString().Replace('"', '""')
        winget complete --word="$Local:word" --commandline "$Local:ast" --position $cursorPosition | ForEach-Object {
            [System.Management.Automation.CompletionResult]::new($_, $_, 'ParameterValue', $_)
        }
  }
}
else
{
  Write-Warning "winget not found on PATH - skipping winget autocomplete and upgrade checks."
}
