# profile.d/90-personal.ps1
# Personal touches from the author (dontrolle). Feel free to delete or replace this file -
# it is not required by anything else in the profile.

Set-Location -Path "c:\src"
